# stage2/04-hardhat/files/hardhat.service
[Unit]
Description=hardhat
After=multi-user.target

[Service]
Type=idle
ExecStart=/root/hardhat.sh
Restart=no

[Install]
WantedBy=multi-user.target